# backup/crontab
# EasyMed database backup - runs every 4 hours
0 */4 * * * . /etc/environment; /backup.sh >> /var/log/backup.log 2>&1
# Empty line required by cron
